# systemd/mplayer-daemon.service
[Unit]
Description=MPlayer daemon

[Service]
Type=simple
ExecStart=/usr/bin/mplayer-daemon

[Install]
WantedBy=multi-user.target
