# udev/99-audio-cd.rules
SUBSYSTEM=="block", KERNEL=="sr0", ENV{ID_CDROM_MEDIA_CD}=="1", ACTION=="change", RUN+="/usr/local/bin/play-cd"
